# perception/supervisord.conf
; Runs the three perception services as separate processes inside one container.
; Each is an independent uvicorn app on its own port — split into its own
; container later by giving each a Dockerfile and dropping supervisor.
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:yolo]
command=uvicorn services.yolo.main:app --host 0.0.0.0 --port 8001
directory=/app/perception
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:yoloseg]
command=uvicorn services.yoloseg.main:app --host 0.0.0.0 --port 8007
directory=/app/perception
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sam3]
command=uvicorn services.sam3.main:app --host 0.0.0.0 --port 8002
directory=/app/perception
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:locateanything]
command=uvicorn services.locateanything.main:app --host 0.0.0.0 --port 8003
directory=/app/perception
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; DDS cloud-detector proxy (T-Rex2 / DINO-X / Grounding-DINO / DINO-XSeek).
; CPU-only network client — no GPU. Needs DDS_API_TOKEN in the environment;
; without it the process still runs and /health responds (loaded=false).
[program:dds]
command=uvicorn services.dds.main:app --host 0.0.0.0 --port 8008
directory=/app/perception
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
